# systemd/discoclaw.service
[Unit]
Description=Discoclaw
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/code/discoclaw
# Keep secrets local; this file should exist on the host (not committed).
EnvironmentFile=%h/code/discoclaw/.env
# Build first: `pnpm build` (this service runs `dist/`).
ExecStart=/usr/bin/node %h/code/discoclaw/dist/index.js
Restart=on-failure
RestartSec=15
StartLimitIntervalSec=600
StartLimitBurst=3

[Install]
WantedBy=default.target
